# playbooks/roles/qat_driver/files/qat_service.service
[Unit]
Description=Intel QuickAssist Technology service

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/etc/init.d/qat_service start
ExecStop=/etc/init.d/qat_service stop

[Install]
WantedBy=multi-user.target
